# scripts/archive/staging-auto-deploy.service
[Unit]
Description=Auto-deploy staging for spheroseg-app
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=cvat
Group=cvat
WorkingDirectory=/home/cvat/spheroseg-app
ExecStart=/home/cvat/spheroseg-app/scripts/auto-deploy-staging.sh
Restart=always
RestartSec=10
StandardOutput=append:/var/log/staging-auto-deploy.log
StandardError=append:/var/log/staging-auto-deploy.log

[Install]
WantedBy=multi-user.target