# scripts/spheroseg-backup.timer
[Unit]
Description=Daily SpheroSeg database backup (pairs with spheroseg-backup.service)
Documentation=https://github.com/michalprusek/cell-segmentation-hub/blob/main/docs/database-backup.md

[Timer]
# Daily at 03:30 local time. PostgreSQL load is lowest overnight; offset
# from common cron times (03:00, 04:00) to avoid pile-up with other host
# maintenance jobs.
OnCalendar=*-*-* 03:30:00

# If the host was off at the scheduled time, run as soon as it boots —
# missing a backup window is more expensive than running late.
Persistent=true

# Spread up to 5 minutes so multiple machines on the same network don't
# hit the database at exactly 03:30:00.
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
